# config/comic-ai-daemon.service
[Unit]
Description=Comic AI Enhanced Daemon - 容錯拓撲和自進化系統
Documentation=https://github.com/anomalyco/opencode
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=codespace
WorkingDirectory=/workspaces/cosmic-ai.uk
ExecStart=/home/codespace/.python/current/bin/python3 /workspaces/cosmic-ai.uk/enhanced_daemon.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=comic-ai-daemon

# Resource limits
MemoryLimit=512M
CPUQuota=80%

# Security
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
